# arcade-feedback.service
[Unit]
Description=Arcade Feedback Server
After=network.target

[Service]
User=hsad
Group=hsad
Environment=HOME=/home/hsad
Environment=PATH=/home/hsad/.local/bin:/usr/local/bin:/usr/bin:/bin
# Allow binding privileged ports (80/443) while running as hsad
AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecStart=/usr/bin/python3 /home/hsad/Arcade/feedback_server.py
WorkingDirectory=/home/hsad/Arcade
Restart=always
RestartSec=5
StandardOutput=append:/home/hsad/Arcade/logs/feedback_server.log
StandardError=append:/home/hsad/Arcade/logs/feedback_server.log

[Install]
WantedBy=multi-user.target
